# packages/bsp/olinuxino/lib/systemd/system/olinuxino-modules.service
[Unit]
Description=Olinuxino modules load serveice
After=systemd-modules-load.service

[Service]
ExecStart=/bin/bash /usr/bin/olinuxino_modules.sh start

[Install]
WantedBy=sysinit.target
