# packages/bsp/olinuxino/lib/systemd/system/olinuxino-bluetooth.service
[Unit]
Description=Realtek H5 bluetooth support
Before=bluetooth.service

[Service]
ExecStart=/bin/bash /usr/bin/olinuxino_bt.sh start
ExecStop=/bin/bash /usr/bin/olinuxino_bt.sh stop

[Install]
WantedBy=multi-user.target
